# files/services/pttest.timer
[Unit]
Description=Run the hotfix instance of Project Tracker

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
